# packages/server/lightpanel.service
[Unit]
Description=NeoPixel Light Panel
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/github/neopixel-light-panel/packages/server
ExecStart=/home/pi/.nvm/versions/node/v14.21.3/bin/node app.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
